# deploy/dashboard.service
[Unit]
Description=Bluesky Collector Dashboard
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bluesky-collector
ExecStart=/home/ubuntu/bluesky-collector/.venv/bin/python dashboard.py --port 5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
