# deploy/prune.timer
[Unit]
Description=Run Bluesky pruner hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
